# monitor/influx-grafana/supervisord.conf
[supervisord]
nodaemon = true

[program:nginx]
command = /usr/sbin/nginx
autorestart = true

[program:influxdb]
command = /opt/influxdb/run.sh
autorestart = true