# monitor/influx-grafana/supervisord.conf
[supervisord]
nodaemon = true

[program:nginx]
command = /usr/sbin/nginx
autorestart = true

[program:influxdb]
command = /usr/local/bin/run_influxdb
autorestart = true